# backup/crontab
0 1 * * * /bin/sh /usr/local/bin/pgsqlbackup.sh 2>&1
0 2 * * * /bin/sh /usr/local/bin/drop_old_dumps.sh 2>&1
