# backup/crontab
0 1 * * * /bin/sh /usr/local/bin/pgsqlbackup.sh 2>&1
